# solar-clock.service
[Unit]
Description=Solar Smart Clock
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/solar-clock

# Load API key from secrets file (not in git)
EnvironmentFile=/etc/solar-clock/secrets

ExecStart=/usr/bin/python3 -m solar_clock
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/dev/fb1 /var/log/solar-clock

[Install]
WantedBy=multi-user.target
