# solar-clock.service
[Unit]
Description=Solar Smart Clock
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/solar-clock
Environment=OPENWEATHER_API_KEY=ac5189e6a0f50737d3145e449c96c5e6
ExecStart=/usr/bin/python3 /home/pi/solar-clock/clock.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
